# magicmirror/rootfs/etc/services.d/magicmirror/run
#!/usr/bin/with-contenv bashio
# ==============================================================================
# Community Hass.io Add-ons: MagicMirror²
# Runs MagicMirror²
# ==============================================================================
# shellcheck disable=SC1091
source /usr/lib/hassio-addons/base.sh

bashio::log.info "Starting MagicMirror² server..."
cd /opt || bashio::exit "Could not change dir."
exec node serveronly
